# tools/com.opencode.dropbox-recovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opencode.dropbox-recovery</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ianhandley/code/github/opencode/tools/auto_recover_dropbox.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Run daily at 2 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/ianhandley/.opencode/logs/dropbox-recovery.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/ianhandley/.opencode/logs/dropbox-recovery.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>AUTO_RECOVER</key>
        <string>true</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
